# monitoring.service
[Unit]
Description=Process monitoring script runner

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitor_test_process.sh
